# cookbooks/node_oracledb/files/oracle-api-service.conf
start on (local-filesystems)
stop on shutdown

script

cd /vagrant/oracle-api
exec forever start index.js

end script
